# roles/silverblue/templates/mount-systemd.automount.j2
[Unit]
Description={{ item.name }}

[Automount]
Where={{ item.mount.where }}

[Install]
WantedBy=multi-user.target